# files.centos7/etc/supervisor.d/crond.conf
[program:crond]
command=/usr/sbin/crond

